# com.bigfundj.nestldown-roster.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bigfundj.nestldown-roster</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/paulburchfield/miniconda3/bin/python3</string>
        <string>/Users/paulburchfield/Documents/projects/dj-availability-checker/nestldown_roster.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/nestldown-roster.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/nestldown-roster.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
